# procfile2supervisor/supervisord.conf
{% for p in processes -%}
[program:{{ app }}]
command                  = {{ p.cmd }}
user                     = {{ user }}
environment              = {% for key,value in p.env.items() -%}
                               {{ key }}={{ value|shellquote }}
                               {%- if not loop.last %},{% endif %}
                           {%- endfor %}
startsecs                = 30
redirect_stderr          = true
autorestart              = true
stopwaitsecs             = 10
directory                = /opt/webapp/%(program_name)s
stdout_logfile           = {{ log }}/%(program_name)s-app-out.log
stdout_logfile_maxbytes  = 10MB
stdout_logfile_backups   = 10
stdout_capture_maxbytes  = 1MB
stderr_logfile           = {{ log }}/%(program_name)s-app-err.log
stderr_logfile_maxbytes  = 10MB
stderr_logfile_backups   = 10
stderr_capture_maxbytes  = 1MB
{%- endfor %}

